# templates/gogs-bkp.timer
[Unit]
Description=Create Gogs Backups

[Timer]
OnCalendar=Wed *-*-* 03:00:00
Unit={{ __gogs_bkp_unit }}

[Install]
WantedBy=timers.target
